# files/rsyslogd
#!/sbin/openrc-run

description="System logging daemon"
supervisor=supervise-daemon

command="/usr/bin/rsyslogd"
command_args="-n"
pidfile="/run/rsyslogd.pid"

depend() {
    need localmount
    after clock bootmisc
    provide logger
}
