# srcpkgs/oracle-cloud-agent/files/etc/init.d/oracle-cloud-agent
#!/sbin/openrc-run

description="Oracle Cloud Infrastructure Agent"
supervisor=supervise-daemon

command="/usr/lib/oracle-cloud-agent/agent"
command_args="-agent-config /etc/oracle-cloud-agent/agent.yml"
supervise_daemon_args="--stdout /var/log/oracle-cloud-agent/agent.log --stderr /var/log/oracle-cloud-agent/agent.log"

depend() {
    need net dhcpcd
    after cloud-final
}

start_pre() {
    # Runtime directories expected by the agent and its plugins
    local dirs="
        /var/log/oracle-cloud-agent
        /var/log/oracle-cloud-agent/plugins
        /var/log/oracle-cloud-agent/plugins/gomon
        /var/log/oracle-cloud-agent/plugins/unifiedmonitoring
        /var/log/oracle-cloud-agent/plugins/oci-osmh
        /var/log/oracle-cloud-agent/plugins/oci-vulnerabilityscan
        /var/log/oracle-cloud-agent/plugins/oci-blockautoconfig
        /var/log/oracle-cloud-agent/plugins/oci-wlp
        /var/log/oracle-cloud-agent/plugins/oci-hpc
        /var/log/oracle-cloud-agent/plugins/oci-hpc/oci-hpc-configure
        /var/log/oracle-cloud-agent/plugins/oci-hpc/oci-hpc-monitoring
        /var/log/oracle-cloud-agent/plugins/oci-hpc/oci-rdma-authentication
        /var/lib/oracle-cloud-agent
        /var/lib/oracle-cloud-agent/tmp
        /var/lib/oracle-cloud-agent/reattach
        /var/lib/oracle-cloud-agent/plugins
        /var/lib/oracle-cloud-agent/plugins/unifiedmonitoring
        /var/lib/oracle-cloud-agent/plugins/oci-osmh
        /var/lib/oracle-cloud-agent/plugins/oci-vulnerabilityscan
        /var/lib/oracle-cloud-agent/plugins/oci-blockautoconfig
        /var/lib/oracle-cloud-agent/plugins/oci-wlp
    "
    for d in $dirs; do
        checkpath -d -m 0755 "$d"
    done

    # bastions plugin requires cap_dac_override to manage session files
    local bastions=/usr/lib/oracle-cloud-agent/plugins/bastions/bastions
    if [ -x "$bastions" ] && command -v setcap >/dev/null 2>&1; then
        setcap cap_dac_override+ep "$bastions" 2>/dev/null || true
    fi
}
